# 61-sddm-plasmawayland.rules
# disable Wayland if modesetting is disabled
RUN+="/usr/bin/rm -f /etc/sddm/hide-wayland-sessions"
IMPORT{cmdline}="nomodeset", RUN+="/usr/bin/touch /etc/sddm/hide-wayland-sessions"
